# venus.cron
* * * * * . /Users/robinsmith/Programming/robin/zsh/projects/dotfiles/functions.zsh; email robin.smith@cloudreach.co.uk "Vancouver Shopping" "Assign to Robin." >> /Users/robinsmith/out.log 2>> /Users/robinsmith/errors.log



